# utils/systemd/borg-backup-sched.service
[Unit]
Description=Scheduled backups with borg

[Service]
ExecStart=/backup/backup-scripts/backup.sh
