# utils/systemd/borg-backup-sched.service
[Unit]
Description=Scheduled backups with borg

[Service]
ExecStart=/backup/borg-scripts/backup.sh
